# rootfs/etc/services.d/mitmproxy/run
#!/usr/bin/execlineb -P
mitmdump -w /config/last.dump
